# service/dbwatcher.service
[Unit]
Description=DB watcher

[Service]
ExecStart=python /home/pi/dbmeter/dbwatcher.py

[Install]
WantedBy=multi-user.target
